# templates/default/pivotal_agent_init.erb
#!/bin/sh
#
# pivotal_agent        init file for starting up the new relic pivotal agent as  daemon
#
# chkconfig: 2345 20 80
# description: Starts and stops the new relic pivotal agent daemon.

# Source function library.
. /etc/rc.d/init.d/functions

PIVOTAL_AGENT=<%= node[:newrelic_monitoring][:pivotal_agent][:name] %> 
PIVOTAL_AGENT_DIR=<%= node[:newrelic_monitoring][:agents_install_dir]%>/<%=node[:newrelic_monitoring][:pivotal_agent][:name]%>
PIVOTAL_AGENT_USER=<%= node[:newrelic_monitoring][:user]%>

# you shouldn't have to change anything below this line

PIVOTAL_EXEC="${PIVOTAL_AGENT_DIR}/${PIVOTAL_AGENT}"

PIDFILE="/var/run/${PIVOTAL_AGENT}.pid"
LOCKFILE="/var/lock/subsys/${PIVOTAL_AGENT}"

start() {
    [ -f $PIVOTAL_EXEC ] || exit 6
    [ -x $PIVOTAL_EXEC ] || exit 5
    echo -n $"Starting $PIVOTAL_AGENT: "
    if [ -z "$USER" -o "$USER" = "root" ]; then
        daemon --user ${PIVOTAL_AGENT_USER} "$PIVOTAL_EXEC"
    else
        daemon $PIVOTAL_EXEC
    fi
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $LOCKFILE
    return $retval
}

stop() {
    echo -n $"Stopping $PIVOTAL_AGENT: "
    killproc -p $PIDFILE $PIVOTAL_AGENT
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $LOCKFILE
    return $retval
}

restart() {
    stop
    start
}

reload() {
    false
}

rh_status() {
    status -p $PIDFILE $name
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}


case "$1" in
    start)
        rh_status_q && exit 0
        $1
        ;;
    stop)
        rh_status_q || exit 0
        $1
        ;;
    restart)
        $1
        ;;
    reload)
        rh_status_q || exit 7
        $1
        ;;
    force-reload)
        force_reload
        ;;
    status)
        rh_status
        ;;
    condrestart|try-restart)
        rh_status_q || exit 0
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart}"
        exit 2
esac
exit $?
